# etc/profile.d/perl-extras_1.0u_linux26-x86-64.sh
[[ -n /home/valin/ssm-domains-base/setup/env_003/perl-extras_1.0u_linux26-x86-64/lib/perl5/site_perl ]] && export PERL5LIB=/home/valin/ssm-domains-base/setup/env_003/perl-extras_1.0u_linux26-x86-64/lib/perl5/site_perl:/home/valin/ssm-domains-base/setup/env_003/perl-extras_1.0u_linux26-x86-64/lib/perl5/site_perl
[[ -z /home/valin/ssm-domains-base/setup/env_003/perl-extras_1.0u_linux26-x86-64/lib/perl5/site_perl ]] && export PERL5LIB=/home/valin/ssm-domains-base/setup/env_003/perl-extras_1.0u_linux26-x86-64/lib/perl5/site_perl
